# crontab
# m h  dom mon dow   command
* * * * * /home/pi/Pumper/uptime-metric.sh

